# openwrt/files/mqtt.lora.init
#!/bin/sh /etc/rc.common

# OpenWRT lora-mqtt init script

START=70

start() {
    echo start
    service_start /usr/bin/mqtt.lora -d
    service_start /usr/bin/mqtt.lora -d
}

stop() {
    echo stop
    service_stop /usr/bin/mqtt.lora
    service_stop /usr/bin/mqtt.lora
}